# screenshot.service
# /etc/systemd/system/screenshot.service
[Unit]
Description=Screenshot Sync Service
After=graphical.target

[Service]
Type=simple
User=zombie
# Setze hier Deine DISPLAY/XAUTHORITY (für X11) bzw. Wayland‐Infos:
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/zombie/.Xauthority
Environment=XDG_RUNTIME_DIR=/run/user/1000
# Bei Wayland könntest Du statt DISPLAY noch setzen:
# Environment=XDG_SESSION_TYPE=wayland
# Environment=WAYLAND_DISPLAY=wayland-0

ExecStart=/home/zombie/zephyr/screenshot.sh
Restart=always

[Install]
WantedBy=multi-user.target
